# rootconf/etc/rc.local
#!/bin/sh

# CUSTOM SYSTEM INIT SCRIPT
# ADD YOUR COMMANDS HERE

